# sqlite-web/rootfs/etc/services.d/sqlite-web/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Runs SQLite Web
# ==============================================================================
declare -a options

options+=(--host 127.0.0.1)
options+=(--no-browser)
options+=(-x "/config/home-assistant_v2.db")

bashio::log.info 'Starting SQLite Web...'

exec sqlite_web "${options[@]}"
