# sqlite-web/rootfs/etc/services.d/sqlite-web/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: SQLite Web
# Runs SQLite Web
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

declare -a options

options+=(-H 0.0.0.0)
options+=(--no-browser)
options+=(-x "/config/$(hass.config.get 'database_path')")

if hass.config.true 'read_only'; then
  options+=(-r)
fi

hass.log.info 'Starting SQLite Web'

exec sqlite_web "${options[@]}"
