# sqlite-web/rootfs/etc/services.d/sqlite-web/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Runs SQLite Web
# ==============================================================================
declare -a options
declare database

database="/config/home-assistant_v2.db"
if bashio::config.has_value 'database'; then
    database="$(bashio::config 'database')"
fi

options+=(--host 127.0.0.1)
options+=(--no-browser)
options+=(-x "${database}")

bashio::log.info 'Starting SQLite Web...'

exec sqlite_web "${options[@]}"
