# sqlite-web/rootfs/etc/services.d/sqlite-web/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Runs SQLite Web
# ==============================================================================
declare -a options
folder_path=/config/
if bashio::config.has_value 'folder_path'; then
    folder_path=$(bashio::config 'folder_path')
fi
options+=(-H 0.0.0.0)
options+=(--no-browser)
options+=(-x "$folder_path$(bashio::config 'database_path')")

if bashio::config.true 'read_only'; then
  options+=(-r)
fi

bashio::log.info 'Starting SQLite Web...'

exec sqlite_web "${options[@]}"
