# sqlite-web/rootfs/etc/services.d/datasette/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: SQLite Web
# Runs Datasette
# ==============================================================================

if bashio::config.false 'datasette'; then
    while true; do
      bashio::log.debug 'Datasette are not enabled...'
      exec sleep 86400
    done
fi

args="/config/$(bashio::config 'database_path')"

bashio::log.debug "Running Datasette with: ${args}"

exec datasette serve "${args}"
