# sqlite-web/rootfs/etc/s6-overlay/s6-rc.d/sqlite-web/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Runs SQLite Web
# ==============================================================================
declare -a options
declare database

database="/homeassistant/home-assistant_v2.db"
if bashio::config.has_value 'database'; then
    database="$(bashio::config 'database')"
fi

options+=(--host 127.0.0.1)
options+=(--no-browser)
options+=(-x "${database}")
options+=(--url-prefix $(bashio::addon.ingress_entry))

bashio::log.info 'Starting SQLite Web...'

exec sqlite_web "${options[@]}"
