# sqlite-web/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: SQLite Web
# Runs the NGINX daemon
# ==============================================================================

# Wait for SQLLite Web to start
bashio::net.wait_for 8080

bashio::log.info "Starting NGINX..."
exec nginx
